# dashboard/launchd/com.minecraft-ai.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.minecraft-ai.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>server.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Volumes/DevDrive/Projects/Minecraft-AI-on-Linux/dashboard</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/myroproductions</string>
    <key>SSH_KEY_PATH</key>
    <string>/Users/myroproductions/.ssh/id_ed25519</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/myroproductions/Library/Logs/minecraft-ai-dashboard/dashboard.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/myroproductions/Library/Logs/minecraft-ai-dashboard/dashboard.err.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
